# deploy/quanthunt.service
# systemd unit for the QuantHunt gunicorn service.
# setup-vm.sh installs this to /etc/systemd/system/quanthunt.service and
# substitutes __DEPLOY_USER__ with the deploy user (opc / ubuntu).
#
# Binds 127.0.0.1 only — the app is never exposed directly; nginx (port 80/443)
# is the sole public entry point and reverse-proxies to it.
[Unit]
Description=QuantHunt screener (gunicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=/opt/quanthunt
Environment=PORT=5000
# Optional secrets/config (e.g. EODHD_API_KEY=...) — the leading '-' means the
# service still starts if the file is absent. Deploys never overwrite it.
EnvironmentFile=-/opt/quanthunt/.env
# One worker, many threads: the universe/price caches are per-process in memory,
# so multiple workers would each re-fetch from NSE and defeat the 6-hour cache.
# Threads (not workers) absorb the burst of parallel calls a single dossier /
# analyser open fires (strategy + timeframes + checklist + fundamentals + …),
# each of which can block for seconds on yfinance/NSE. 8 threads saturated and
# hung the whole site; 32 gives ample headroom for several concurrent users.
# --graceful-timeout recycles a worker that a stuck request wedges.
ExecStart=/opt/quanthunt/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 1 --threads 32 --timeout 120 --graceful-timeout 30 wsgi:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
